# service.d/00xvfb/run
#!/bin/bash
exec /usr/bin/Xvfb :1 -screen 0 ${XRFBRESX}x${XRFBRESY}x${XRFBDEPTH}