# service.d/01wm/run
#!/bin/bash
#
#
sv start 00xvfb || exit 1
export HOME="/home/user"
export DISPLAY=":1"
export USER="user"
exec sudo -E -u $USER ratpoison --display :1
