# service/02x11vnc/run
#!/bin/bash
exec /usr/bin/x11vnc -display :1 -nopw -wait 5 -forever
