# service.d/01chrome/run
#!/bin/bash
HOME="/home/chrome"
DISPLAY=":1"
USER="chrome"
exec sudo -E -u $USER /opt/google/chrome/chrome --user-data-dir  --window-position=0,0 --window-size=1024,768 --force-device-scale-factor=1
